# packages/common/postinst
#!/bin/sh

# Create user and group
if ! id -u sudobot >/dev/null 2>&1; then
    echo "Creating user 'sudobot'..."
    useradd --system --home /var/sudobot --shell /bin/false sudobot
fi

if ! getent group sudobot >/dev/null 2>&1; then
    echo "Creating group 'sudobot'..."
    groupadd --system sudobot
    usermod -aG sudobot sudobot
fi

# Create directories
if [ ! -d /var/sudobot ]; then
    mkdir -p /var/sudobot
    mkdir -p /var/sudobot/config
    mkdir -p /var/sudobot/tmp

    echo "{}" > /var/sudobot/config/config.json
    echo "{}" > /var/sudobot/config/system.json
    chmod 644 /var/sudobot/config/config.json
    chmod 644 /var/sudobot/config/system.json
fi

# Create directories
if [ ! -d /var/log/sudobot ]; then
    mkdir -p /var/log/sudobot
fi

chown -R sudobot:sudobot /var/log/sudobot
chown -R sudobot:sudobot /var/sudobot
chown -R sudobot:sudobot /usr/share/sudobot
chown -R sudobot:sudobot /etc/sudobot
chmod 600 /etc/sudobot/env

# Install dependencies
cd /usr/share/sudobot || {
    echo "Failed to change directory to /usr/share/sudobot." >&2
}

if [ ! -d /usr/share/sudobot/node_modules ]; then
    echo "Installing dependencies in /usr/share/sudobot..."
fi

su sudobot -c "cd /usr/share/sudobot; PATH=\"$PATH\" HOME=/var/sudobot npm install --omit=dev" -s /bin/sh || {
    echo "Failed to setup sudobot." >&2
}

# Enable systemd service
if [ -f /lib/systemd/system/sudobot.service ] && command -v systemctl >/dev/null 2>&1; then
    if command -v deb-systemd-helper >/dev/null 2>&1; then
        deb-systemd-helper enable sudobot.service || true
    fi
fi

# Create SystemV init script
if command -v update-rc.d >/dev/null 2>&1; then
    cat << 'EOF' > /etc/init.d/sudobot
#!/bin/sh
### BEGIN INIT INFO
# Provides:          sudobot
# Required-Start:    $local_fs $network
# Required-Stop:     $local_fs $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start sudobot daemon
# Description:       Start the sudobot daemon
### END INIT INFO

exec /usr/bin/sudobotctl "$@"
EOF
    chmod 755 /etc/init.d/sudobot
    update-rc.d sudobot defaults || true
fi
