# packages/common/postinst
#!/bin/sh

# Create user and group
if ! id -u sudobot >/dev/null 2>&1; then
    echo "Creating user 'sudobot'..."
    useradd --system --home /var/sudobot --shell /bin/false sudobot
fi

if ! getent group sudobot >/dev/null 2>&1; then
    echo "Creating group 'sudobot'..."
    groupadd --system sudobot
    usermod -aG sudobot sudobot
fi

# Create directories
if [ ! -d /var/sudobot ]; then
    mkdir -p /var/sudobot
fi

chown -R sudobot:sudobot /var/sudobot
chown -R sudobot:sudobot /usr/share/sudobot
chown -R sudobot:sudobot /etc/sudobot
chmod 600 /etc/sudobot/env

# Install dependencies
cd /usr/share/sudobot || {
    echo "Failed to change directory to /usr/share/sudobot." >&2
}

pm=""

if [ -n "$SUDO_USER" ] && [ -d "/home/$SUDO_USER" ]; then
    if [ -d "/home/$SUDO_USER/.bun" ]; then
        if [ -z "$BUN_INSTALL" ]; then
            BUN_INSTALL="/home/$SUDO_USER/.bun"
        fi
        
        PATH="$BUN_INSTALL/bin:$PATH"
    fi
fi

if command -v pnpm >/dev/null 2>&1; then
    pm="pnpm"
elif command -v bun >/dev/null 2>&1; then
    pm="bun"
elif command -v npm >/dev/null 2>&1; then
    pm="npm"
elif command -v yarn >/dev/null 2>&1; then
    pm="yarn"
else
    echo "No package manager found. Please install pnpm, bun, npm, or yarn." >&2
    exit 1
fi

su sudobot -c "PATH=\"$PATH\" BUN_INSTALL=/var/sudobot $pm install" -s /bin/sh || {   
    echo "Failed to install dependencies with $pm." >&2
}

# Enable systemd service
if [ -f /lib/systemd/system/sudobot.service ]; then
    if command -v deb-systemd-helper >/dev/null 2>&1; then
        deb-systemd-helper enable sudobot
    fi
fi