# packages/common/prerm
#!/bin/sh

if command -v deb-systemd-helper >/dev/null 2>&1; then
    deb-systemd-helper stop sudobot.service
    deb-systemd-helper disable sudobot.service
fi

if [ -d /var/sudobot ] && [ -n "$(ls -A /var/sudobot)" ]; then
    echo "Note: /var/sudobot is not empty. It will not be removed."
fi

chown -R root:root /var/sudobot
chown -R root:root /etc/sudobot

rm -rf /usr/share/sudobot/node_modules

if id -u sudobot >/dev/null 2>&1; then
    userdel sudobot
fi

if getent group sudobot >/dev/null 2>&1; then
    groupdel sudobot
fi